# templates/zot.service.epp
[Unit]
Description=OCI Distribution Registry
Documentation=https://zotregistry.io/
After=network.target auditd.service local-fs.target

[Service]
Type=simple
ExecStart=<%= $binary %> serve <%= $config_file %>
Restart=on-failure
User=<%= $user %>
Group=<%= $group %>
LimitNOFILE=<%= $limit_nofile %>
<%- if $memory_high { -%>
MemoryHigh=<%= $memory_high %>
<%- } -%>
<%- if $memory_max { -%>
MemoryMax=<%= $memory_max %>
<%- } -%>

[Install]
WantedBy=multi-user.target
